# .other/disk-writer/etc/sv/disk-writer/run
#!/bin/sh

exec /usr/bin/node /usr/share/disk-writer/app 2>&1
